# dist/autoaurrepo.service
[Unit]
Description=autoaurrepo

[Service]
Type=simple
ExecStart=autoaurrepo