# dist/autoaurrepo.service
[Unit]
Description=autoaurrepo

[Service]
Type=oneshot
ExecStart=autoaurrepo